# deploy/cacp.service
[Unit]
Description=CACP cloud room server
After=network.target

[Service]
Type=simple
User=cacp
Group=cacp
WorkingDirectory=/opt/cacp
EnvironmentFile=/etc/cacp/cacp.env
ExecStart=/usr/bin/corepack pnpm --filter @cacp/server start
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/var/lib/cacp

[Install]
WantedBy=multi-user.target
